# dnsgraph.cron
# 
*/1 * * * *    /usr/sbin/rndc stats
*/10 * * * *   /var/lib/named/bin/dnsgraph.pl > /var/lib/named/named-graph.log
